# systemd/mochad.service
[Unit]
Description=Start mochad service

[Service]
Type=forking
ExecStart=/usr/local/bin/mochad

[Install]
WantedBy=multi-user.target
