# udev/mochad.service
[Unit]
Description=Mochad X10 TCP gateway
After=syslog.target
After=sockets.target

[Service]
ExecStart=/usr/local/bin/mochad
